# Times/build/files/times.service
# systemd service configuration

[Unit]
Description=Times Service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/python /opt/times/times.py
Restart=always

[Install]
WantedBy=multi-user.target
