# SensorHub/sensorhub.service
# systemd service configuration

[Unit]
Description=SensorHub Service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/python /usr/local/sonar/sensorhub.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
